# config/supervisor.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisor.log
logfile_maxbytes=10MB
loglevel=info
pidfile=/var/run/supervisord.pid

[program:php-fpm]
command=/usr/sbin/php-fpm7.0 --nodaemonize -c /etc/php/7.0/fpm/
autostart=true
autorestart=unexpected
stderr_logfile=/var/log/supervisor/php-fpm.err.log
stderr_logfile_maxbytes=10MB
stdout_logfile_maxbytes=10MB
stdout_logfile=/var/log/supervisor/php-fpm.log

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=unexpected
stderr_logfile=/var/log/supervisor/nginx.err.log
stderr_logfile_maxbytes=10MB
stdout_logfile=/var/log/supervisor/php-fpm.log
stdout_logfile_maxbytes=10MB
